# spec/fixtures/samba-ad/supervisord.conf
[inet_http_server]
port = 9001
username = bolt
password = bolt

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl = http://127.0.0.1:9001

[program:samba]
command   = /usr/sbin/samba -i

[supervisord]
logfile   = /tmp/supervisord.log   ; supervisord log file
loglevel  = error                  ; info, debug, warn, trace
nodaemon  = true
user      = root
pidfile   = /tmp/supervisord.pid   ; pidfile location
directory = /tmp
